# supervisor.conf
[program:sample_project]
command = /webapps/sample_project/bin/gunicorn_start
user = sample_project
stdout_logfile = /webapps/sample_project/logs/gunicorn_supervisor.log
redirect_stderr = true